# ops/maintenance/line-bridge-healthcheck.timer
[Unit]
Description=Run LINE Discord Bridge health checks every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=line-bridge-healthcheck.service

[Install]
WantedBy=timers.target
